# catalog/crontab
12 3 * * * python /code/manage.py refresh --all