# catalog/crontab
12 3 * * * python /code/manage.py refresh --all
5 */1 * * * python /code/manage.py timeout --ago 6
12 20 * * * python /code/manage.py truncate --ago 30